# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/local/bin/gpiod.lua &
/usr/local/bin/monitord.lua &

exit 0
